# cron/crontab
37 4 * * * root /usr/local/sbin/renew-certs
